# config/caps2esc
#!/sbin/openrc-run

name="caps2esc"
command="/usr/sbin/caps2escd"
